# scripts/mysterymixclub-api-prod.service
# MysteryMixClub FastAPI backend — production.
# Install on the Droplet AS mysterymixclub-api.service (same unit name as
# staging uses on its own box — one env per Droplet, so there's no collision):
#   sudo cp scripts/mysterymixclub-api-prod.service /etc/systemd/system/mysterymixclub-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mysterymixclub-api
[Unit]
Description=MysteryMixClub FastAPI (production)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=mysterymixclub
Group=mysterymixclub
WorkingDirectory=/home/mysterymixclub/app/backend
EnvironmentFile=/etc/mysterymixclub/prod.env
ExecStart=/home/mysterymixclub/app/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
